# rtl8852au-dkms-git/40-rtl8852au-usb-modeswitch.rules
ATTR{idVendor}=="2357", ATTR{idProduct}=="0141", RUN+="usb_modeswitch '/%k'"
